# roles/docker/templates/prune.j2.timer
[Timer]
OnCalendar=Hourly

[Install]
WantedBy=timers.target